# ttrpg-pi.service
[Unit]
Description=TTRPG Pi Server
After=network.target graphical.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/ttrpg-pi
ExecStart=/usr/bin/python3 /home/pi/ttrpg-pi/ttrpg_pi.py
Restart=on-failure
RestartSec=10
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/pi/.Xauthority"

[Install]
WantedBy=graphical.target
